# statistics-gathering/client/etc/cron.hourly/0nep52-publisher
#!/bin/bash
    /usr/local/sbin/nep52-publisher <vm-type> clock
